# python_code/templogger/jessie_temp_logger/templogger.service
[Unit]
Description=TempLogger

[Service]
Type=simple
ExecStart=/home/pi/templogger.py
