# python_code/weather_screen/jessie_weather_screen/weather_screen.service
[Unit]
Description=weather_screen

[Service]
Type=simple
ExecStart=/home/pi/jessie_weather_screen.py
