# python_code/shutdown_button_code/pullup_shutdown_jessie_pi/pishutdown.service
[Service]
ExecStart=/usr/bin/python /home/pi/shutdown/pishutdown.py
WorkingDirectory=/home/pi/shutdown/
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=pishutdown
User=root
Group=root

[Install]
WantedBy=multi-user.target
